# deploy/renderflow-backup.service
[Unit]
Description=RenderFlow daily backup (pg_dump + projects tar)
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /opt/renderflow/deploy/backup.sh
